# slimbook/DEBIAN/postinst
#!/bin/bash
systemctl daemon-reload
systemctl enable slimbook.service
systemctl restart slimbook.service
